# systemd/surveil-value-directory.timer
[Unit]
Description=Run MarketPulseWire ValueList investment-bank research collector at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=surveil-value-directory.service

[Install]
WantedBy=timers.target
