# systemd/surveil-news-collector-shadow.service
[Unit]
Description=MarketPulseWire news media collector shadow batch
After=network-online.target surveil-db-init.service surveil-proxy.service
Wants=network-online.target surveil-proxy.service

[Service]
Type=oneshot
User=surveil
WorkingDirectory=/opt/surveil
EnvironmentFile=-/opt/surveil/.env
EnvironmentFile=-/opt/surveil/proxy.env
ExecStart=/opt/surveil/.venv/bin/python /opt/surveil/scripts/news_collector.py --write-report --limit 5 --direct-shadow
StandardOutput=append:/opt/surveil/logs/news-collector-shadow.log
StandardError=append:/opt/surveil/logs/news-collector-shadow.err.log
