# systemd/surveil-company-disclosures.timer
[Unit]
Description=Run Surveil company disclosures at 08:00 and 20:00

[Timer]
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 20:00:00
Persistent=true
Unit=surveil-company-disclosures.service

[Install]
WantedBy=timers.target
