# systemd/surveil-news-collector-shadow.timer
[Unit]
Description=Run MarketPulseWire news media collector shadow batch every 10 minutes

[Timer]
OnActiveSec=3min
OnBootSec=9min
OnUnitActiveSec=10min
Persistent=true
Unit=surveil-news-collector-shadow.service

[Install]
WantedBy=timers.target
